# services/ffmpeg_cam.service.template
[Unit]
Description=FFmpeg Stream für Kamera {{CAM_NUM}}
After=network.target

[Service]
ExecStart={{SCRIPT_PATH}}
Restart=always
User={{USERNAME}}

[Install]
WantedBy=multi-user.target